# sp-elasticsearch.conf
[program:elasticsearch]
command=/home/gitpod/elasticsearch-7.9.3/bin/elasticsearch -d -p /home/gitpod/elasticsearch-7.9.3/pid -Ediscovery.type=single-node
process_name=%(program_name)s
priority=1001
autostart=true
startretries=3
autorestart=true
user=gitpod
