# sp-rabbitmq.conf
[program:rabbitmq]
command=/etc/init.d/rabbitmq-server start
process_name=%(program_name)s
priority=1001
autostart=true
startretries=3
autorestart=true
user=root
