# sp-redis.conf
[program:redis]
command=/usr/bin/redis-server
process_name=%(program_name)s
priority=1001
autostart=true
startretries=3
autorestart=true
user=root
